# Containers/vpn/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dnsmasq]
command=dnsmasq --no-daemon --log-facility=-
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dnsmasq.log
stderr_logfile=/var/log/supervisor/dnsmasq.log

[program:wireguard]
command=/opt/vpn-sdk/vpn-manager.sh start
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/wireguard.log
stderr_logfile=/var/log/supervisor/wireguard.log