# Containers/web3-dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ganache]
command=ganache --host 0.0.0.0 --port 8545 --deterministic --accounts 10 --mnemonic "test test test test test test test test test test test junk"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ganache.log
stderr_logfile=/var/log/supervisor/ganache.log

[program:web3-workspace]
command=/bin/bash -c "cd /workspace && echo 'Web3 workspace ready at /workspace'"
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/workspace.log
stderr_logfile=/var/log/supervisor/workspace.log